# roles/monitoring/sw-pve-monitor/templates/proxmox-mon.timer.tmpl
[Unit]
Description=Run proxmox data export every 5 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=300
AccuracySec=1ms

[Install]
WantedBy=timers.target 